# root/etc/services.d/cups/run
#!/usr/bin/with-contenv bash

exec s6-setuidgid abc cupsd -f -c /config/cupsd.conf -s /config/cups-files.conf
